# aria2/linux/aria2_update_tracker.timer
[Unit]
Description=Update Aria2c BitTorrent Tracker List Automatically

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
